# dnscrypt/services.conf
[program:dnscrypt]
command = dnscrypt-proxy -a 127.0.0.1:2053 -R opendns

[program:dnsmasq]
command = dnsmasq -d
