# hass/arm/ibeacon.timer
[Unit]
Description=Hourly restart ibeacon

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
