# webkit/supervisor/conf.d/restart.conf
[program:restart]
command = /usr/local/bin/restart.sh
