# webkit/supervisor/conf.d/webkit.conf
[program:webkit]
command = /usr/local/bin/phantomjs server.js %(process_num)d
numprocs_start = 8080
numprocs = 10
process_name = %(program_name)s-%(process_num)d
directory = /webkit/
