# ggr/data/supervisor.d/confd.ini
[program:confd]
command = confd -watch -backend file -file /dev/null
redirect_stderr = true

[program:null]
command = sh -c 'while true; do date > /dev/null; sleep 60; done'
redirect_stderr = true
