# ggr/data/supervisor.d/confd.ini
[program:confd]
command = confd -watch -backend file -file /dev/null
