# ggr/data/supervisor.d/ggr.ini
[program:ggr]
command = ggr -listen :4444 -users /etc/grid-router/users.htpasswd -quotaDir /etc/grid-router/quota -guests-allowed
redirect_stderr = true
