# ggr/data/supervisor.d/null.ini
[program:null]
command = sh -c 'while true; do date > /dev/null; sleep 60; done'
