# dnscrypt/proxy/services.conf
[program:dnscrypt]
priority = 100
command = dnscrypt-proxy -a 127.0.0.1:2053 -R opendns

[program:pdnsd]
priority = 200
command = pdnsd -c /etc/pdnsd.conf
