# webkit/etc/supervisor/conf.d/haproxy.conf
[program:haproxy]
command = /usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg -p /var/run/haproxy.pid
