# webkit/etc/supervisor/conf.d/monit.conf
[program:monit]
command = /usr/bin/monit -I -c /etc/monit/monitrc
