# webkit/etc/supervisor/conf.d/monit.conf
[program:monit]
command = /usr/local/bin/monit -I -c /etc/monitrc
