# proxyhub/services.ini
[program:dnscrypt]
priority = 100
command = dnscrypt-proxy -a 127.0.0.1:2053 -R opendns

[program:pdnsd]
priority = 200
command = pdnsd -c /etc/pdnsd.conf

[program:shadowsocks]
priority = 300
command = bash -c 'sslocal -c <(jq "select(.local_port==%(process_num)d)" /etc/shadowsocks.json)'
numprocs_start = 1981
numprocs = 5
process_name = %(process_num)d

[program:haproxy]
priority = 400
command = haproxy -f /etc/haproxy/haproxy.cfg

[program:polipo]
priority = 500
command = /usr/bin/polipo -c /etc/polipo/config

[program:nginx]
priority = 600
command = nginx -g "daemon off; error_log stderr;"

