# dnscrypt/proxy/dnscrypt-proxy/DEBIAN/postinst
#!/bin/bash

ldconfig
systemctl start dnscrypt-proxy
systemctl enable dnscrypt-proxy
