# dnscrypt/proxy/dnscrypt-proxy/DEBIAN/prerm
#!/bin/bash

systemctl disable dnscrypt-proxy
systemctl stop dnscrypt-proxy
