# ggr-ui/data/supervisor.d/ggr-ui.ini
[program:ggr-ui]
command = ggr-ui -listen 8888 -limit 10 -quota-dir /etc/grid-router/quota
redirect_stderr = true
