# renderjs/supervisor/conf.d/renderjs.conf
[program:renderjs]
command = /usr/local/bin/phantomjs /usr/local/bin/render.js %(process_num)d
numprocs_start = 8080
numprocs = 10
process_name = %(program_name)s-%(process_num)d
